# supervisor-windows.conf
[program:test]
command = ./go_test_api.exe
